# scripts/installer/linux/config/upstart-centos.conf
start on started sshd
stop on runlevel [016]
respawn
exec changeme
