# scripts/installer/linux/config/upstart-ubuntu.conf
start on started ssh
stop on runlevel [016]
respawn
exec changeme
